# s6-rc.d/sdrplay-api/run
#!/command/execlineb -P
/opt/bin/sdrplay_apiService