# s6-rc.d/soapysdr-remote/run
#!/bin/sh
export LD_LIBRARY_PATH=/opt/lib
# 需要先检测再开启
/opt/bin/SoapySDRUtil --probe || exit
/opt/bin/SoapySDRServer --bind=$SOPAY_REMOTE_BIND || exit